# users/Profpatsch/whatcd-resolver/services/reverse-proxy/run
#!/usr/bin/env execlineb
caddy run --config ./caddyfile
